# kasparService.conf
description "KASPAR Control Serivce"
author  "Nathan Burke <n.burke@natbur.com>"

start on runlevel [2345]
stop on runlevel [016]

chdir /opt/RobotControl
exec /opt/RobotControl/Web/site.py
respawn